# host_files/flanagan/quadlets/valkey.container
[Container]
Image=cgr.dev/chainguard/valkey
Network=mail.network
Volume=valkey.volume:/data:U,Z
AutoUpdate=registry
[Install]
WantedBy=multi-user.target
